# tags/2.1.44/C++/RailDriver/systemd/raildriverd-hidraw_AT.service
[Unit]

[Service]
Type=forking
GuessMainPID=false
ExecStart=/bin/sh -c "exec /lib/udev/raildriverd-hidraw add %I"
