# branches/2.1.39/C++/RailDriver-libusb/90-raildriver.rules.in
# Raildriver 
ATTRS{idVendor}=="05f3", ATTRS{idProduct}=="00d2", MODE="0666", RUN+="raildriverd $env{ACTION} $env{DEVPATH}"

