# branches/2.1.29.3/C++/RailDriver/90-raildriver.rules.in
# Raildriver 
BUS=="usb", SYSFS{idVendor}=="05f3", SYSFS{idProduct}=="00d2", ACTION=="add", MODE="0666", RUN+="raildriverd $env{ACTION} $env{DEVPATH}"

